# source/net/netbase-1.17-1/debian.conffiles
/etc/init.d/netbase
/etc/gateways
/etc/hosts.allow
/etc/hosts.deny
/etc/rpc
